# .docker/supervisor/supervisord.d/01-httpd.ini
[program:httpd]
process_name=%(program_name)s
command=httpd -DFOREGROUND
autostart=true
autorestart=true
